# deploy/continuityos-backup.service
[Unit]
Description=ContinuityOS data backup

[Service]
Type=oneshot
ExecStart=/home/scott/ai-workspace/repos/continuityos/scripts/backup_data.sh
